# apps/runner/install/s6/runner/run
#!/bin/sh
exec 2>&1
cd /app
exec node dist/index.js
